# com.nucleusbridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nucleusbridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/richardgraham/.pyenv/versions/3.10.12/bin/python3</string>
        <string>/Users/richardgraham/nucleusbridge/main.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/richardgraham/nucleusbridge/bridge.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/richardgraham/nucleusbridge/bridge.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/richardgraham/nucleusbridge</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
